# deploy_galera/templates/etc/systemd/system/mysqld.service.d/override.conf.j2
[Unit]
After=time-sync.target
 
[Service]
LimitNOFILE= "{{ LimitNOFILES }}"